# gunicorn.conf
[program:gunicorn]
command=gunicorn codesandbox:app
directory=/path/to/project
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
